# build/config/blueprint_readonly.conf
# Upstart script at /etc/init/blueprint.conf

start on runlevel [2345] and stopped cloud-final
stop on runlevel [016]

respawn
respawn limit 10 5

script
  # use su <user> because that makes /etc/environment variables available to subprocess
  exec su root -c '/opt/science/blueprint/bin/run_blueprint.sh  -bind=":8000" -readonly=true'
  emit blueprint_readonly_running
end script

