# build/config/schema_suggestor.conf
# Upstart script at /etc/init/schema_suggestor.conf

start on runlevel [2345] and stopped cloud-final
stop on runlevel [016]

respawn
respawn limit 10 5

script
  # use su <user> because that makes /etc/environment variables available to subprocess
  exec su root -c '/opt/science/blueprint/bin/run_schema_suggestor.sh'
  emit schema_suggestor_running
end script

